# docker/conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autorestart=true
autostart=true