# docker/conf.d/lxpanel.conf
[program:lxpanel]
command=/usr/bin/lxpanel --profile LXDE
environment=DISPLAY=":1",HOME="/root",USER="root"
autorestart=true
user=root
priority=15