# docker/conf.d/lxpanel.conf
[program:lxpanel]
command=/usr/bin/lxpanel --profile LXDE-pi
environment=DISPLAY=:1,HOME="/root",USER="root"
autorestart=true
autostart=true
user=root
priority=1
directory=/root