# docker/conf.d/lxpanel.conf
[program:lxpanel]
autorestart=true
command=/usr/bin/lxpanel --profile LXDE
environment=DISPLAY=":1",HOME="/root",USER="root"
priority=15
user=root