# docker/conf.d/novnc.conf
[program:novnc]
autorestart=true
command=/root/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 8080
priority=25