# docker/conf.d/novnc.conf
[program:novnc]
priority=25
command=/root/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 8080
autorestart=true