# docker/conf.d/openbox.conf
[program:openbox]
autorestart=true
command=openbox
environment=DISPLAY=":1",HOME="/root",USER="root"
priority=15