# docker/conf.d/openbox.conf
[program:openbox]
command=openbox
autorestart=true
environment=DISPLAY=":1",HOME="/root",USER="root"
priority=15