# docker/conf.d/pcmanfm.conf
[program:pcmanfm]
priority=15
command=/usr/bin/pcmanfm --desktop --profile LXDE
user=root
environment=DISPLAY=":1",HOME="/root",USER="root"