# docker/conf.d/pcmanfm.conf
[program:pcmanfm]
command=/usr/bin/pcmanfm --desktop --profile LXDE
environment=DISPLAY=":1",HOME="/root",USER="root"
priority=15
user=root