# docker/conf.d/x11vnc.conf
[program:x11vnc]
autorestart=true
command=x11vnc -display :1 -forever -shared -xkb -repeat -capslock
priority=20