# docker/conf.d/x11vnc.conf
[program:x11vnc]
priority=20
command=x11vnc -display :1 -forever -shared -xkb -repeat -capslock
autorestart=true