# docker/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :1 -forever -shared -xkb 
autorestart=true
autostart=true