# docker/conf.d/xterm.conf
[program:xterm]
command=xterm
autorestart=true
autostart=true